# initcpio/rc.local
#!/bin/sh
# Add local users
cat > /etc/passwd <<EOF
root::0:0::/root:/bin/sh
nobody:x:99:99:Unprivileged User:/dev/null:/bin/false
EOF

# Add some name servers
cat > /etc/resolv.conf <<EOF
nameserver 1.1.1.1
nameserver 1.0.0.1
EOF

# Mount devpts
mkdir -p /dev/pts
echo "devpts           /dev/pts              devpts  gid=5,mode=620  0 0" >> /etc/fstab
mount /dev/pts
